# bot/monitor_cpu.service
[Unit]
Description=Monitor CPU Usage Service
After=network.target

[Service]
ExecStart=/usr/local/bin/monitor_cpu_usage.sh
Restart=on-failure

[Install]
WantedBy=multi-user.target